# supervisord.conf
[program:auto_recover]
command=/health/auto_recover.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/auto_recover.err.log
stdout_logfile=/var/log/supervisor/auto_recover.out.log
