# health/health/metrics/metrics/supervisord.conf
[program:metrics_server]
command=/metrics/metrics_server.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/metrics.err.log
stdout_logfile=/var/log/supervisor/metrics.out.log
